# airootfs/etc/skel/.config/autostart/kitest-welcome.desktop
[Desktop Entry]
Type=Application
Name=Kitest Welcome
Exec=/usr/local/bin/kitest-welcome
X-GNOME-Autostart-enabled=true
