# .tmp/kitest-work/x86_64/airootfs/etc/xdg/autostart/kitest-live-first-run.desktop
[Desktop Entry]
Type=Application
Name=Kitest live first-run
Comment=Apply Catppuccin on first login + Flathub extras (best-effort)
Exec=/usr/local/bin/kitest-live-first-run.sh
X-KDE-autostart-phase=2
