# airootfs/etc/xdg/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=Install Kitest OS
Comment=Launch Calamares installer (starts after Plasma login)
Exec=/usr/local/bin/kitest-calamares-safe
TryExec=/usr/local/bin/kitest-calamares-safe
Icon=calamares
Terminal=false
StartupNotify=true
Categories=Qt;System;
Hidden=false
X-GNOME-Autostart-enabled=true
