# airootfs/etc/xdg/autostart/kitest-live-first-run.desktop
[Desktop Entry]
Type=Application
Name=Kitest live first-run
Comment=Flathub + Kitest Flatpak bundles (best-effort)
Exec=/usr/local/bin/kitest-live-first-run.sh
TryExec=/usr/local/bin/kitest-live-first-run.sh
OnlyShowIn=KDE;
X-KDE-autostart-phase=2
